# debian/prerm
#!/bin/sh
set -eu

case "$1" in
  remove|deconfigure)
    systemctl stop heikodiscopi.service || true
    systemctl disable heikodiscopi.service || true
    ;;
esac

exit 0
